# .config/systemd/user/backup-reminder.service
[Unit]
Description=Remind us to make a backup if we did not create/update one recently

[Service]
Type=oneshot
ExecStart=%h/bin/backup --reminder
