# .config/systemd/user/grep-exclude-cache-invalidation.service
[Unit]
Description=Invalidate grep(1) exclude cached files whenever .cvsignore is modified

[Service]
Type=simple
ExecStart=/usr/bin/rm -rf %C/grep
