# .config/systemd/user/grep-exclude-cache-invalidation.path
[Unit]
Description=Invalidate grep(1) exclude cached files whenever .cvsignore is modified

[Path]
PathChanged=%h/.cvsignore

[Install]
WantedBy=paths.target
