# .config/systemd/user/update-system-reminder.service
[Unit]
Description=Remind us to update the system if we did not recently

[Service]
Type=oneshot
ExecStart=%h/bin/update-system --reminder
