# .config/autostart/user-session.desktop
[Desktop Entry]
Type=Application
Name=user-session
Exec=sh -c "$HOME/bin/user-session"
Terminal=false
Comment=
RunHook=0
