# .config/etc/systemd/system/lightdm.service.d/99-local.conf
[Service]

# We reset `Restart`  to `no` to prevent the session  from being restarted after
# pressing `M-SysRq-e`.
Restart=no
