# .config/etc/systemd/system/getty@.service.d/99-activate-numlock.conf
[Service]
# `setleds(1)`: set the keyboard leds for the console
# `-D`: also change default setting
# `+num`: set NumLock flag
ExecStartPre=/usr/bin/setleds -D +num
StandardInput=file:/dev/%I
